# app/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:PhSpring]
command=sh /home/PhSpring/bin/start_server.sh start -d